# ansible/roles/geoserver/templates/geoserver.service.j2
[Unit]
Description=GeoServer {{ geoserver_role | upper }} ({{ inventory_hostname }})
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
User=geoserver
Group=geoserver
# Garante acesso ao socket Docker independente da ordem de criação do usuário/grupo
SupplementaryGroups=docker
WorkingDirectory=/opt/geoserver
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=geoserver-{{ geoserver_role }}
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
